# docker-compose/docker/host-cron/crontab
0 0 * * 0 root sh /send_mail_weekly.sh
* * * * * /send_mail_weekly.sh
# Single container cron example